# etc/crontab/reboot
#!/bin/bash

# /bin/bash -c "echo 255 > /sys/devices/pwm-fan/target_pwm"
